# crontab
0 */12 * * * {PATH_TO_PYTHON3.9<=} main.py
